# supervisor/kafka.conf
[program:kafka]
command=/usr/bin/start-kafka.sh
loglevel=debug
stderr_logfile=/var/log/supervisor/kafka-stderr.log
autostart=true
autorestart=true
