# supervisor/zookeeper.conf
[program:zookeeper]
command=/usr/bin/start-zookeeper.sh
loglevel=debug
stderr_logfile=/var/log/supervisor/zookeeper-stderr.log
autostart=true
autorestart=true
